# examples/units/user/swaylid.service
[Unit]
Description=Swaylid

[Service]
ExecStart=~/.cargo/bin/swaylid

[Install]
WantedBy=graphical.target
